# nvidia-control.service
[Unit]
Description=NVIDIA GPU Control Panel
After=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/bin/nvidia-control-gui
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target